# lazylibrarian/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# Called by s6 when nginx exits. Non-zero exit code means unexpected crash.
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.error "nginx exited unexpectedly with code ${1}"
fi
