# lazylibrarian/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

port=5299

# Wait for LazyLibrarian to become available
bashio::net.wait_for "$port" localhost 900

echo ""
if bashio::config.equals "connection_mode" "noingress"; then
    bashio::log.info "Ingress disabled. Connect at http://homeassistant:$(bashio::addon.port 5299)"
else
    bashio::log.info "Ingress enabled. Open from the HA sidebar or http://homeassistant:$(bashio::addon.port 5299)/lazylibrarian"
fi
echo ""

bashio::log.info "Starting nginx..."
exec nginx
